# cmd/agent/internal/phases/host/assets/unbounded-agent-daemon.service
[Unit]
Description=Unbounded Agent Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/unbounded-agent daemon --endpoint={{.Endpoint}}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
