# cmd/agent/internal/phases/host/assets/unbounded-agent-daemon.service
# Copyright (c) Microsoft Corporation.
# Licensed under the MIT License.

[Unit]
Description=Unbounded Agent Daemon
After=network-online.target machines.target
Wants=network-online.target machines.target

[Service]
Type=simple
ExecStart=/usr/local/bin/unbounded-agent-current daemon
Restart=always
RestartSec=10
StartLimitIntervalSec=0
OnFailure=unbounded-agent-daemon-recovery.service

[Install]
WantedBy=multi-user.target
