# cmd/agent/internal/phases/host/assets/unbounded-agent-daemon.service
[Unit]
Description=Unbounded Agent Daemon
After=network-online.target machines.target
Wants=network-online.target machines.target

[Service]
Type=simple
Environment=UNBOUNDED_AGENT_CONFIG_FILE=/etc/unbounded-agent/config.json
ExecStart=/usr/local/bin/unbounded-agent daemon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
